# LaunchAgents/com.backbone.granola-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.backbone.granola-sync</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>.backbone/granola.js</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/elliotstorey/backbone-v9</string>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>0</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/elliotstorey/backbone-v9/.backbone/granola-sync.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/elliotstorey/backbone-v9/.backbone/granola-sync.log</string>

	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
